# docker/s6-overlay/resources/services.d/worker/finish
#!/command/with-contenv bash

# Couleurs ANSI
YELLOW='\033[0;33m'
RED='\033[0;31m'
GREEN='\033[0;32m'
NC='\033[0m' # No Color

# Vérifier si le worker est désactivé
if [ "$ENABLE_WORKER" != "true" ]; then
    echo -e "${GREEN}✅ Worker service désactivé - pas de redémarrage automatique${NC}"
    exit 125  # Code spécial s6-overlay pour ne pas redémarrer
fi

echo -e "${YELLOW}🔄 Laravel Queue Worker arrêté, redémarrage automatique...${NC}"

# Vérifier le code de sortie
if [ $1 -ne 0 ]; then
    echo -e "${RED}⚠️  Worker arrêté avec le code d'erreur: $1${NC}"
fi

# Code de sortie 0 pour redémarrer automatiquement
exit 0

############################################################################